# files/munin/ejabberd_registrations.cron
#!/bin/bash
if [ -x "/usr/sbin/ejabberdctl" ]; then
  EJCTL="/usr/sbin/ejabberdctl"
else
  EJCTL=$(ls -1 /opt/ejabberd-*/bin/ejabberdctl | head -n 1)
fi

for host in `cat /etc/munin/plugin-conf.d/ejabberd_registrations.conf  | grep env.vhosts | sed 's/env.vhosts //'`; do
  stat_file="/var/lib/munin/plugin-state/ejabberd_registrations_${host}"
  $EJCTL registered-users $host > $stat_file
  chgrp munin $stat_file
  chmod o-rwx $stat_file
done
